# bacrypt/bacrypt.service
[Unit]
Description=Backup and encrypt files

[Service]
Type=oneshot
ExecStart=/home/bacrypt/timed.sh
